# lib/systemd/system/wpa_supplicant.service.d/snap.conf
[Unit]
ConditionPathExists=!/etc/systemd/system/snap.wpa-suplicant.wpa.service
