# config/lib/systemd/system/snapd.generate-network-conf.service
[Unit]
Description=Generate default network config

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/snapd/generate-network-conf
